# cron.d
* * * * * root /app/hello.sh >> /var/log/cron.log